# iot-portal.service
[Unit]
Description=IoT Portal Flask Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/iot-portal
ExecStart=/usr/bin/python3 /opt/iot-portal/app.py
Restart=always
RestartSec=3
Environment=PYTHONDONTWRITEBYTECODE=1

[Install]
WantedBy=multi-user.target
